# deploy/systemd/icicle-lending.service
[Unit]
Description=Icicle Lending Liquidation-Risk Engine
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__HOME__/icicle

# Config (CLICKHOUSE_PASSWORD, ICICLE_ARCHIVE_RPC). On setup.sh-provisioned hosts
# this file holds them. The leading - makes it optional so hosts that inject the
# same vars via systemd drop-ins (Environment=) work without this file present.
EnvironmentFile=-/etc/icicle/lending.env

ExecStart=__HOME__/icicle/icicle lending --chain 43114 --metrics-port 9092

Restart=always
RestartSec=10
LimitNOFILE=524288

# Mostly an in-memory account/state cache plus batched RPC reads. Modest caps.
MemoryHigh=2G
MemoryMax=4G

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
