# deploy/icicle-api.service
[Unit]
Description=Icicle API Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/icicle
ExecStart=/root/icicle/icicle api --port 8080
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
